# alerts-service/scripts/alerts-service.service
[Unit]
Description=Agent8080 Whale Alerts Service
After=network.target

[Service]
Type=simple
User=claw
WorkingDirectory=/home/claw/.openclaw/workspace/projects/agent-claw-wallet/alerts-service
Environment=NODE_ENV=production
EnvironmentFile=/home/claw/.openclaw/workspace/projects/agent-claw-wallet/alerts-service/.env
ExecStart=/usr/bin/node src/index.js
Restart=always
RestartSec=10
StandardOutput=append:/home/claw/.openclaw/workspace/projects/agent-claw-wallet/alerts-service/logs/alerts-service.log
StandardError=append:/home/claw/.openclaw/workspace/projects/agent-claw-wallet/alerts-service/logs/alerts-service.log

[Install]
WantedBy=multi-user.target